# crontab
*/30 * * * * /home/romi-devops/check_disk.sh
